# packaging/playbooks/roles/abbot/templates/abbot.socket.j2
[Socket]
ListenStream={{ abbot_port }}

[Install]
WantedBy=sockets.target
