# packaging/playbooks/roles/holyorders/templates/holyorders.timer.j2
[Unit]
Description="Run Holy Orders hourly"

[Timer]
OnBootSec=15min
OnUnitActiveSec={{ holyorders_frequency }}

[Install]
WantedBy=timers.target
