# packaging/abbot.socket
[Socket]
ListenStream=8888

[Install]
WantedBy=sockets.target
