# packaging/playbooks/abbott.socket.j2
[Socket]
ListenStream={{ abbott_port }}

[Install]
WantedBy=sockets.target
